# cron/crontab
*/5 * * * * /home/candyy/cron/low_battery_notify.sh
