# cron/crontab
*/5 * * * * /home/candyy/stuff/cron/low_battery_notify.sh
